# sdcard-wifi-ssh/initrd_extracted/etc/systemd/system/keybow.service
[Unit]
Description=Keybow Daemon
After=usb-gadget.service

[Service]
ExecStart=/start.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
